# send_ip/templates/get_net_iface.timer.j2
[Unit]
Description=Run get_net_iface.py at 08:30 Mon-Fri

[Timer]
OnCalendar=Mon..Fri 08:30
Persistent=true

[Install]
WantedBy=timers.target
